# source/root/root/.config/autostart/compton compositor.desktop
[Desktop Entry]
Exec=compton
Name=compton compositor
Type=Application
Version=1.0
X-LxQt-Need-Tray=true
