# source/root/etc/systemd/system/dnscrypt-proxy.socket.d/HSF.conf
[Socket]
# own this socket even when network is inactive
FreeBind=true
MaxConnections=256
#NoDelay=false
#DeferAcceptSec=0
# additions
ListenStream=127.1.0.1:53
ListenDatagram=127.1.0.1:53
