# source/root/etc/profile.d/beep-disable.sh
# doesn't work
#set bell-style none

# the least squeaky beep
setterm --bfreq 25 2> /dev/null

# shut it
setterm --blength 0 2> /dev/null
