# source/root/etc/profile.d/zzz-jack.csh
setenv SDL_AUDIODRIVER alsa
setenv AUDIODRIVER alsa
setenv QEMU_AUDIO_DRV alsa
setenv ALSA_CONFIG_PATH /etc/asound-jack.conf
