# source/root/etc/profile.d/zzz-jack.csh
unsetenv SDL_AUDIODRIVER
unsetenv QEMU_AUDIO_DRV
setenv AUDIODRIVER alsa
setenv ALSA_CONFIG_PATH=/etc/alsa-jack.conf
