# source/root/etc/profile.d/zzz-jack.csh
unsetenv SDL_AUDIODRIVER
setenv AUDIODRIVER alsa
