# source/root/etc/profile.d/zzz-jack.sh
unalias sox
export SDL_AUDIODRIVER=alsa
alias timidity='timidity -Oj'
export AUDIODRIVER=alsa
export QEMU_AUDIO_DRV=alsa
export ALSA_CONFIG_PATH=/etc/asound-jack.conf
