# data/common-userfiles/.config/autostart/akregator.desktop
[Desktop Entry]
Exec=akregator --hide-mainwindow
Name=Akregator
Type=Application
Version=1.0
X-LXQt-Need-Tray=false
Hidden=true
