# data/common-userfiles/.config/autostart/keepnote.desktop
[Desktop Entry]
Exec=keepnote
Name=keepnote
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
