# data/common-userfiles/.config/autostart/LinPhone.desktop
[Desktop Entry]
Exec=linphone --iconified
Hidden=true
Name=LinPhone
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
